# systemd/crafty-server-watcher.service
[Unit]
Description=Crafty Server Watcher — auto-hibernate and wake Minecraft servers
Documentation=https://github.com/your-repo/crafty-server-watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/opt/crafty-server-watcher/venv/bin/python -m crafty_server_watcher --config /etc/crafty-server-watcher/config.yaml
WorkingDirectory=/opt/crafty-server-watcher

# ── Credentials ──────────────────────────────────────────────────
EnvironmentFile=-/etc/crafty-server-watcher/env

# ── Service identity ─────────────────────────────────────────────
User=crafty-watcher
Group=crafty-watcher

# ── Sandboxing ───────────────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
ReadWritePaths=/var/log/crafty-server-watcher

# ── Capabilities (only needed if binding to port < 1024) ───────
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# ── Restart policy ──────────────────────────────────────────────
Restart=on-failure
RestartSec=10

# ── Resource limits ─────────────────────────────────────────────
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
